# debian-src/DEBIAN/postinst
# #!/bin/sh
# set -e

# # Only try user stuff if we have a session bus
# i f [ -n "$XDG_RUNTIME_DIR" ] && [ -n "$DBUS_SESSION_BUS_ADDRESS" ]; then
#     i f c ommand -v systemctl >/dev/null 2>&1; then
#         s ystemctl --user daemon-reload || true
#         systemctl --user enable --now hotspot-daemon.service || true
#     fi

#     if c ommand -v gnome-extensions >/dev/null 2>&1; then
#         gnome-extensions enable hotspot-indicator@hotspot || true
#     fi
# fi

# exit 0